# config/picoclaw-funnel.service
[Unit]
Description=Tailscale Funnel for PicoClaw
After=tailscaled.service

[Service]
ExecStart=/usr/bin/tailscale funnel 18790
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
