# systemd/user/picoclaw-watchdog.service
[Unit]
Description=PicoClaw Ops Watchdog (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash -lc 'if [ -f "%h/.picoclaw/.env" ]; then set -a; . "%h/.picoclaw/.env"; set +a; fi; exec "%h/.local/share/picoclaw/scripts/ops_watchdog.sh" once'
WorkingDirectory=%h
StandardOutput=append:%h/.picoclaw/logs/ops-watchdog.log
StandardError=append:%h/.picoclaw/logs/ops-watchdog.log

[Install]
WantedBy=default.target
